# stow/systemd/.config/systemd/user/ss-local.service
[Unit]
Description=Shadowsocks Local Proxy
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ss-local -v -c /etc/shadowsocks-libev/config.json
Restart=on-failure
RestartSec=5

# Security hardening (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target

